# crontab
* * * * * root bash /app/script.sh >> /var/log/cron.log 2>&1
# new cron job
